# configure-wakeup.service
[Unit]
After=local-fs.target

[Service]
ExecStart=/usr/local/bin/configure_wakeup.sh

[Install]
WantedBy=default.target
